# deb/src/deb/control/postinst
#!/bin/sh
/usr/bin/systemctl daemon-reload
/usr/bin/systemctl enable phone-home-detector --now
