# deb/src/deb/control/prerm
#!/bin/sh
/usr/bin/systemctl stop phone-home-detector
/usr/bin/systemctl disable phone-home-detector --now
